# templates/default/redis@.service.erb
[Unit]
Description=Redis In-Memory Data Store
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
ExecStart=<%= @bin_path %>/redis-server /etc/redis/redis.conf
ExecStop=<%= @bin_path %>/redis-cli shutdown
Restart=always

[Install]
WantedBy=multi-user.target
